# scripts/systemd/dcx-backend.service
# DCX Agent 백엔드 서비스 (FastAPI + Gunicorn)
# 설치: sudo cp scripts/systemd/dcx-backend.service /etc/systemd/system/
#       sudo systemctl daemon-reload
#       sudo systemctl enable --now dcx-backend

[Unit]
Description=DCX Agent Backend (FastAPI)
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/dcx_agent/backend
Environment="PATH=/home/ubuntu/dcx_agent/backend/venv/bin:/usr/bin"
EnvironmentFile=/home/ubuntu/dcx_agent/.env

ExecStart=/home/ubuntu/dcx_agent/backend/venv/bin/gunicorn app.main:app \
    --worker-class uvicorn.workers.UvicornWorker \
    --workers 2 \
    --bind 127.0.0.1:8000 \
    --timeout 300 \
    --graceful-timeout 30 \
    --access-logfile - \
    --error-logfile -

Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dcx-backend

[Install]
WantedBy=multi-user.target
